# scripts/com.prasad.jobtracker.aggregator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prasad.jobtracker.aggregator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/prasadkanade/Documents/Prasad Kanade/Job Hunt Tracker/scripts/cron_runner.sh</string>
        <string>aggregator</string>
    </array>

    <!-- Run at 8 AM, 3 PM, 9 PM every day -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <!-- Catch-up: runs missed jobs when Mac wakes from sleep -->
    <key>StandardOutPath</key>
    <string>/Users/prasadkanade/Documents/Prasad Kanade/Job Hunt Tracker/.local/launchd_aggregator.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/prasadkanade/Documents/Prasad Kanade/Job Hunt Tracker/.local/launchd_aggregator_err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/prasadkanade/Documents/Prasad Kanade/Job Hunt Tracker</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/prasadkanade</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
